# sys-power/power-profiles-daemon/files/power-profiles-daemon.initd
#!/sbin/openrc-run
# Copyright 2022 Calculate Linux by jksf
# Distributed under the terms of the GNU General Public License v2

command="/usr/libexec/power-profiles-daemon"
command_args="${PPDARGS}"
supervisor="supervise-daemon"
command_args_foreground=""

depend() {
	use clock logger
	need dbus
}
